# infra/lint/boss-conservation-invariants.timer
[Unit]
Description=Boss — conservation-invariant sweep (hourly)
Documentation=https://github.com/algedonic-dev/boss/blob/main/docs/design/correctness-protocol.md

[Timer]
# Hourly is a good cadence: long enough that the sweep is
# cheap (a handful of indexed queries) and short enough that
# any drift surfaces inside one operating shift. A future
# bulletin-emit hook will fan failures onto the Operations
# page; today the journal is the only sink.
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
